# config/systemd/system/a7z-ssh-host-keys.service
[Unit]
Description=Generate SSH host keys for the A7Z image
ConditionPathExists=!/etc/ssh/ssh_host_ed25519_key
Before=ssh.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A

[Install]
WantedBy=multi-user.target
